# deploy/systemd/legacy-crossplay-relay.service
[Unit]
Description=Legacy Console Crossplay Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=legacy-relay
Group=legacy-relay
WorkingDirectory=/opt/legacy-crossplay-relay
EnvironmentFile=/etc/legacy-crossplay-relay.env
ExecStart=/usr/bin/dotnet /opt/legacy-crossplay-relay/LegacyCrossplayRelay.dll
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/log/legacy-crossplay-relay
LockPersonality=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
